# config.service
[Unit]
Description=Network FenPing configuration

[Service]
ExecStart=/.boot

[Install]
WantedBy=multi-user.target
